# linux/crontab
# run on reboot
crontab -e
# insert
@reboot <absolutePath>/<file>.sh

# make executable
chmod +x <absolutePath>/<file>.sh

